# deploy/systemd/angel-bot-daemon.service
[Unit]
Description=Angel Bot Runtime Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__APP_USER__
Group=__APP_USER__
WorkingDirectory=__APP_DIR__
Environment=PYTHONUNBUFFERED=1
ExecStart=__APP_DIR__/.venv/bin/python __APP_DIR__/bot_daemon.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
