# deploy/systemd/angel-dashboard.service
[Unit]
Description=Angel Bot Streamlit Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__APP_USER__
Group=__APP_USER__
WorkingDirectory=__APP_DIR__
Environment=PYTHONUNBUFFERED=1
ExecStart=__APP_DIR__/.venv/bin/streamlit run __APP_DIR__/dashboard.py --server.address 0.0.0.0 --server.port 8501
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
